# services/notmuch/run
#!/bin/sh
exec env                         \
    HOME="${HOME:-/home/noctua}" \
    DISPLAY="${DISPLAY:-:0}"     \
    ./notmuch
