# services/notmuch/run
#!/bin/sh
exec env                         \
    HOME="${HOME:-/home/noctua}" \
    PATH="/usr/local/bin:$PATH"  \
    DISPLAY="${DISPLAY:-:0}"     \
    ./notmuch
