# root/etc/sv/runsvuserdir/run
#!/bin/sh
export USER=ninewise
export HOME=/home/ninewise
. $HOME/.profile

groups="$(id -Gn ninewise | tr ' ' ':')"
svdir="$HOME/service"

exec chpst -u "$USER:$groups" runsvdir "$svdir"
