# services/x/run
#!/bin/sh
cd "$HOME"

. "$HOME"/.profile

exec xinit 2>&1
