# services/x/run
#!/bin/sh
HOME=/home/ninewise

cd "$HOME"

. "$HOME"/.profile

# requires allowed_users = anybody (man Xorg.wrap) in /etc/X11/Xorg.wrapper
exec xinit >> "$HOME/xlog" 2>&1
