# services/x/run
#!/bin/sh
cd "$HOME"

. "$HOME"/.yash_profile

exec xinit 2>&1
