# root/etc/sv/userservices/run
#!/bin/sh
exec chpst -u $(groups noctua | tr -s ' ' ':') runsvdir /home/noctua/services 2>&1
