# config/systemd/user/mbsync.timer
[Unit]
Description=Run iSync every 5 minutes

[Timer]
OnBootSec=10
OnUnitActiveSec=5min

[Install]
WantedBy=default.target

