# config/systemd/user/mbsync.service
[Unit]
Description=iSYnc IMAP Synchronizer

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -a

