# config/systemd/user/urxvtd.service
[Unit]
Description=Rxvt Unicode Daemon

[Service]
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=default.target
